# debian/tcpconsole.postinst
#! /bin/sh

systemctl enable tcpconsole

systemctl start tcpconsole
